# tftpd-manager/DEBIAN/postinst
#!/bin/bash
set -e

APP_CONF="/etc/tftpd-manager.conf"

if [ ! -f "$APP_CONF" ]; then
    echo 'TFTP_DIR="/srv/tftp"' > "$APP_CONF"
    chmod 644 "$APP_CONF"
fi

if systemctl is-active --quiet tftpd-hpa; then
    systemctl stop tftpd-hpa >/dev/null 2>&1 || true
fi

echo "Переключение в ручной режим системной службы tftpd-hpa..."
systemctl disable tftpd-hpa >/dev/null 2>&1 || true

chmod +x /usr/local/bin/tftpd-manager

exit 0