# tftpd-manager/DEBIAN/postrm
#!/bin/bash
set -e

APP_CONF="/etc/tftpd-manager.conf"

restore_system_service() {
    if systemctl list-unit-files | grep -q "^tftpd-hpa.service"; then
        echo "Восстановление системной службы tftpd-hpa..."
        systemctl enable tftpd-hpa >/dev/null 2>&1 || true
    fi
}

case "$1" in
    remove)
        restore_system_service
        ;;

    purge)
        if [ -f "$APP_CONF" ]; then
            rm -f "$APP_CONF"
        fi
        restore_system_service
        ;;

    upgrade|failed-upgrade|abort-install|abort-upgrade|disappear)
        ;;

    *)
        exit 0
        ;;
esac

exit 0